# install/com.onvault.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.onvault.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/onvaultd</string>
    </array>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/onvaultd.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/onvaultd.out</string>
</dict>
</plist>
